# devel/etc/catkin/profile.d/25.turtlebot-stage.sh
/home/mech-user/enshu_ws/devel/.private/turtlebot_stage/etc/catkin/profile.d/25.turtlebot-stage.sh